# bot/scripts/algosoft-bot.timer
[Unit]
Description=Start AlgoSoft bot at 08:00 IST (02:30 UTC) Mon-Fri
Documentation=https://github.com/ssrajpal2001/Option_Selling_May_2026
Requires=algosoft-bot.service

[Timer]
# 08:00 AM IST = 02:30 UTC  (IST is UTC+5:30)
OnCalendar=Mon..Fri *-*-* 02:30:00 UTC

# If the server was off at 08:00 AM, fire immediately on next boot
Persistent=true

# A 30 s randomised delay prevents thundering-herd if multiple services start
RandomizedDelaySec=30

Unit=algosoft-bot.service

[Install]
WantedBy=timers.target
